# docker-server-source0400-WebServer/docker-server-source0100-Tomcat/docker-data/tomcatd.conf
[program:tomcatd]
command = /opt/webserver/tomcat/bin/catalina.sh run
