# docker-server-source0700-ReverseProxy/docker-server-source0100-nginx/supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:cron]
command=/usr/sbin/cron -f
